# files/site.conf
[program:site]
command=uwsgi --ini /dmoj/files/uwsgi.ini
directory=/dmoj/site
stopsignal=QUIT
stdout_logfile=/dmoj/logs/site.log
stderr_logfile=/dmoj/logs/site.log
